# s6-overlay/s6-rc.d/redis/run
#!/command/execlineb -P
redis-server --save "" --appendonly no
