# s6-overlay/s6-rc.d/livekit/run
#!/bin/sh
exec 2>&1

# Generate livekit config from environment
LIVEKIT_API_KEY="${LIVEKIT_API_KEY:-devkey}"
LIVEKIT_API_SECRET="${LIVEKIT_API_SECRET:-secret}"

cat > /tmp/livekit.yaml <<EOF
port: 7880
bind_addresses:
  - "0.0.0.0"
keys:
  ${LIVEKIT_API_KEY}: ${LIVEKIT_API_SECRET}
rtc:
  tcp_port: 7881
  udp_port: 7882
  use_external_ip: true
  use_ice_lite: true
redis:
  address: localhost:6379
room:
  auto_create: true
  empty_timeout: 300
  max_participants: 50
logging:
  level: info
turn:
  enabled: false
EOF

exec /usr/local/bin/livekit-server --config /tmp/livekit.yaml
